# sys-auth/python-validity/files/python-validity.initd-r0
#!/sbin/openrc-run
# python-validity OpenRC init script
# Author: Matteo Salonia (@saloniamatteo)

name="python-validity"
description="validity fingerprint sensor support daemon"

command="/usr/lib/python-validity/dbus-service"
command_args=""
command_user="root:root"

pidfile="/var/run/python-validity.pid"
start_stop_daemon_args="-u root -b -m -p ${pidfile}"

depend() {
	need open-fprintd
}
